# trader.service
[Unit]
Description=Raspberry Pi Trading Bot
After=network.target

[Service]
Type=simple
User=sebachowa
WorkingDirectory=/home/sebachowa/code/pi-trader
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/sebachowa/code/pi-trader"
ExecStart=/home/sebachowa/code/pi-trader/venv/bin/python /home/sebachowa/code/pi-trader/run.py --log-level INFO
Restart=always
RestartSec=30

# Resource limits for Raspberry Pi
MemoryLimit=512M
CPUQuota=80%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/sebachowa/code/pi-trader/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trader-pi

[Install]
WantedBy=multi-user.target